# rootfs/etc/supervisor.d/tomcat.ini

[program:tomcat]
command        = /bin/sh -e /opt/tomcat/bin/catalina.sh run
stdout_logfile = /tmp/%(program_name)s.log
stderr_logfile = /tmp/%(program_name)s.log
autorestart    = true

